# INT26-6/systemd/disk_monitor.service
[Unit]
Description=disk monitor every 5min

[Service]
Type=oneshot
ExecStart=/usr/local/bin/disk_monitor.sh