# INT26-6/systemd/ram_monitor.service
[Unit]
Description=ram monitor every 5min

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ram_monitor.sh